# files/logrotate.timer
[Unit]
Description=Hourly Log Rotation

[Timer]
OnCalendar=hourly
AccuracySec=1h
Persistent=true
